# scripts/shareserial-server.service
[Unit]
Description=ShareSerial Server - Remote Serial Port Sharing
After=network.target
Documentation=https://github.com/semidrive/shareserial

[Service]
Type=simple
User=root
Group=dialout
ExecStart=/usr/local/bin/shareserial-server --config /etc/shareserial/server.yaml
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target